# autostart/mockup-php-launch-temp.desktop
[Desktop Entry]
Type=Application
Name=mockup-php-launch-temp
Exec=php