# systemd/admin-interface.service
[Unit]
Description=Admin Interface

[Service]
ExecStart=/bin/bash /usr/sbin/admin-interface.sh

[Install]
WantedBy=multi-user.target
